# deploy/systemd/clarinet-worker@.service
[Unit]
Description=Clarinet Pipeline Worker (%i)
After=network.target rabbitmq-server.service
Wants=rabbitmq-server.service

[Service]
Type=simple
User=clarinet
Group=clarinet
WorkingDirectory=/opt/clarinet
ExecStart=/opt/clarinet/venv/bin/clarinet worker --queues %i --workers 2
Restart=on-failure
RestartSec=5

# Environment
EnvironmentFile=-/opt/clarinet/env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/clarinet /var/log/clarinet /opt/clarinet
PrivateTmp=true

[Install]
WantedBy=multi-user.target
