# deploy/systemd/clarinet-worker@.service
[Unit]
Description=Clarinet Pipeline Worker (%i)
After=network.target rabbitmq-server.service
Wants=rabbitmq-server.service
# In a multi-VM topology /var/lib/clarinet/data is an NFS mount from the stand;
# block startup until it is mounted so task output never lands in a shadowed
# local dir after a reboot. No-op on single-VM (resolves to the root mount).
RequiresMountsFor=/var/lib/clarinet/data

[Service]
Type=simple
User=clarinet
Group=clarinet
WorkingDirectory=/opt/clarinet
ExecStart=/opt/clarinet/venv/bin/clarinet worker --queues %i --workers 2
Restart=on-failure
RestartSec=5

# Environment
EnvironmentFile=-/opt/clarinet/env
Environment=HOME=/opt/clarinet

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/clarinet /var/log/clarinet /opt/clarinet
PrivateTmp=true

[Install]
WantedBy=multi-user.target
